# supervisord.conf
[program:auth_worker]
command=/home/jeager/microservice/auth_service/venv/bin/celery -A auth_service worker -l INFO
stdout_logfile=/tmp/celery.log
stderr_logfile=/tmp/celery.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
